# ops/systemd/invoice-ocr.service
[Unit]
Description=Invoice OCR docker compose stack
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=true
WorkingDirectory=/opt/invoice-ocr
EnvironmentFile=/opt/invoice-ocr/.env
Environment=ENV=__ENV__

# Refresh secrets, then bring the stack up.  Using --remove-orphans so an
# image whose service was renamed gets cleaned up cleanly.
ExecStartPre=/bin/bash /opt/invoice-ocr/ops/pull-secrets.sh
ExecStart=/usr/bin/docker compose -f /opt/invoice-ocr/docker-compose.yml up -d --remove-orphans
ExecStop=/usr/bin/docker compose -f /opt/invoice-ocr/docker-compose.yml down

# Restart docker compose if it ever goes weird; back off so we don't loop hot
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
